# dev.jamiewhite.contextrecall.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Context Recall Launch Agent
        Starts the daemon automatically on user login.

        This file is a template; /Users/USER/ is a placeholder.
        Run "make install" or "./scripts/install.sh" to install with
        the correct paths substituted for your user account.

        Manual installation:
          cp dev.jamiewhite.contextrecall.agent.plist ~/Library/LaunchAgents/
          launchctl load ~/Library/LaunchAgents/dev.jamiewhite.contextrecall.agent.plist

        Removal:
          launchctl unload ~/Library/LaunchAgents/dev.jamiewhite.contextrecall.agent.plist
          rm ~/Library/LaunchAgents/dev.jamiewhite.contextrecall.agent.plist
    -->
    <key>Label</key>
    <string>dev.jamiewhite.contextrecall.agent</string>

    <key>ProgramArguments</key>
    <array>
        <!--
            Production (PyInstaller binary):
              <string>/Applications/Context Recall.app/Contents/Resources/resources/context-recall-daemon/context-recall-daemon</string>

            Development (Python venv):
              <string>/path/to/context-recall/.venv/bin/python3</string>
              <string>-m</string>
              <string>src.main</string>
        -->
        <string>/Applications/Context Recall.app/Contents/Resources/resources/context-recall-daemon/context-recall-daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Applications/Context Recall.app/Contents/Resources/resources/context-recall-daemon</string>

    <!-- Start automatically on login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if the process exits unexpectedly. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Throttle restarts to avoid rapid-fire loops on error. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Log to ~/Library/Logs for easier debugging. -->
    <key>StandardOutPath</key>
    <string>/Users/USER/Library/Logs/Context Recall/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USER/Library/Logs/Context Recall/stderr.log</string>
</dict>
</plist>
